# debian/postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]; then
    echo "configure haproxy systemd wrapper"
	systemctl enable haproxy
fi
